# application/physical/saltstack/srv/salt/storage/mnt-sda2.mount.sample
[Unit]
Description = Western Digital Pi Drive 314GB

[Mount]
What = /dev/sda2
Where = /mnt/sda2
Type = ext4

[Install]
WantedBy = multi-user.target
