# application/physical/saltstack/srv/salt/storage/srv-http.mount.sample
[Unit]
Description = Western Digital Pi Drive 314GB

[Mount]
What = /mnt/sda2/http
Where = /srv/http
Type = ext4
Options = bind

[Install]
WantedBy = multi-user.target
