# application/physical/saltstack/srv/salt/storage/srv-http.automount.sample
[Unit]
Description = Western Digital Pi Drive 314GB

[Automount]
Where = /srv/http

[Install]
WantedBy = multi-user.target
